# target/linux/ramips/mt7620/APP03/base-files/etc/crontabs/root
30 3 * * * reboot
*/2 * * * * welldone.sh
*/1 * * * * check_online.sh
